# devices/marina-bluetooth-monitor.service
[Unit]
Description=Marina Bluetooth Device Monitor
After=network.target bluetooth.target
Wants=bluetooth.target

[Service]
Type=simple
User=adminx
Group=adminx
WorkingDirectory=/home/adminx/Marina/devices
ExecStart=/usr/bin/python3 /home/adminx/Marina/devices/bluetooth_monitor.py --interval 20
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/adminx/Marina

# Permissions for Bluetooth access
SupplementaryGroups=bluetooth

[Install]
WantedBy=multi-user.target
